# com.cal-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cal-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>VENV_PATH/bin/calsync</string>
        <string>--config</string>
        <string>PROJECT_PATH/config.yaml</string>
        <string>--state</string>
        <string>PROJECT_PATH/state.json</string>
    </array>
    <key>WorkingDirectory</key>
    <string>PROJECT_PATH</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>LOG_PATH/cal-sync-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>LOG_PATH/cal-sync-stderr.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
